# configs/profile.d/welcome.sh
# Lab Space 
# A collection of pre-configured docker images for quick testing of code.
# Github: https://www.github.com/awesomelewis2007/lab_space
# Licence: GNU General Public License v3.0
# By: Lewis Evans

echo "To build and run this project type make"
echo "To edit this project you can use nano, vim or vi"
echo "make edit      <-- Opens the main file in nano"
echo "make edit_vim  <-- Opens the src/ directory in vim"
echo "make edit_vi   <-- Opens the main file in vi"
echo "Github: https://github.com/awesomelewis2007/lab_space"
echo "If you need more help type labspace_help"
